# recipes-rover/rover-app/files/rover-app.service
[Unit]
Description = Rover APP
After = network-online.target

[Service]
ExecStart = /usr/bin/roverapp

[Install]
WantedBy = multi-user.target
